# izpbx-asterisk/rootfs/etc/supervisord.d/dnsmasq.ini
[program:dnsmasq]
command=/usr/sbin/dnsmasq -k --log-facility=-
user=root
group=root
autostart=false
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_events_enabled = true
stderr_events_enabled = true
