# izpbx-asterisk/rootfs/etc/supervisord.d/izpbx.ini
[program:izpbx]
command=startizpbx
autostart=false
autorestart=true
startretries=3
startsecs=1
user=root
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_events_enabled = true
stderr_events_enabled = true

[eventlistener:izpbx-exit]
command=bash -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
process_name=izpbx
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL
